# server/fty-discovery-ng.service.in
[Unit]
Description=FTY discovery service

[Service]
Type=simple
ExecStart=${INSTALL_DIR}/${NAME} --config=${CONFIG_DIR}/discovery.conf --daemon
User=root
ExecReload=/bin/kill -HUP $MAINPID
